# scripts/bluetooth.service.d/robot-factory.conf
# robot_factory: shrink BlueZ 5.53 attack surface for BLE peripheral mode.
# Keep the daemon focused on GAP + custom GATT server only. Disable classic
# profiles and GATT client profile plugins that may probe a connected phone.
[Service]
ExecStartPre=
ExecStopPost=
ExecStartPre=/bin/bash /etc/xinit/xbt_46212
ExecStart=
ExecStart=/usr/lib/bluetooth/bluetoothd --noplugin=a2dp,avrcp,audio,sap,input,hog,hostname,policy,health,network,autopair,scanparam,deviceinfo,battery,sixaxis,midi
ExecStopPost=/bin/bash /etc/xinit/xbtex
